# modules/supervisor/templates/supervisor.conf.erb
[program:<%= app_name %>]
command=<%= command %>
directory=<%= directory %>
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=True